# deploy/systemd/siap-worker.service
[Unit]
Description=SIAP queue worker (importaciones y trabajos en cola)
After=network.target mysql.service

[Service]
User=www-data
Group=www-data
Restart=always
RestartSec=5
WorkingDirectory=/var/www/html/siap-puerto-boyaca
ExecStart=/usr/bin/php /var/www/html/siap-puerto-boyaca/artisan queue:work --sleep=3 --tries=3 --max-time=3600
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
